# sw/test_tools/cdk/files/sam-theme-custom/files/root/etc/uci-defaults/luci-theme-sam
#!/bin/sh
uci batch <<-EOF
	set luci.themes.sam=/luci-static/sam
        commit luci
EOF
	
